# deploy/systemd/mag-newspaper-api.service
[Unit]
Description=Newspaper Intelligence API Backend
After=network.target

[Service]
Type=simple
User=magms2596
Group=magms2596
WorkingDirectory=/home/mag.mstatilitechnologies.com/public_html
Environment=PATH=/home/mag.mstatilitechnologies.com/.venv/bin
EnvironmentFile=/home/mag.mstatilitechnologies.com/.env
ExecStart=/home/mag.mstatilitechnologies.com/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8007
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/mag.mstatilitechnologies.com/logs/api.log
StandardError=append:/home/mag.mstatilitechnologies.com/logs/api-error.log
SyslogIdentifier=mag-newspaper-api

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/mag.mstatilitechnologies.com/storage /home/mag.mstatilitechnologies.com/logs

[Install]
WantedBy=multi-user.target