# deploy/systemd/mag-newspaper-archive.timer
[Unit]
Description=Run Newspaper Intelligence PDF Archive Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
